# udev/98-eventlircd4irserver.rules
SUBSYSTEM=="tty", DRIVERS=="ftdi_sio", ATTRS{interface}=="IRTrans USB", ACTION=="add", RUN+="/sbin/initctl emit --no-wait start-irserver"