# etc/init/irexec.conf
description     "IREXEC Daemon"
author          "Steffen Barszus <steffenbpunkt@gmail.com>"

start on started eventlircd
stop on stopping eventlircd
respawn

exec /usr/bin/irexec /etc/lirc/lircrc
